# etc/profile.d/ant.sh
ANT_HOME=/opt/knlservices/Java/apache-ant-1.9.6
export ANT_HOME
PATH=$ANT_HOME/bin:$PATH
export PATH
