# etc/profile.d/java.sh
JAVA_HOME=/opt/knlservices/Java/jdk1.8.0_60
export JAVA_HOME

PATH=$JAVA_HOME/bin:$PATH

export PATH
