# etc/profile.d/jenkins.sh
JENKINS_HOME=/opt/knlservices/Java/jenkins/jenkins-etc

export JENKINS_HOME
