# etc/profile.d/maven.sh
MAVEN_HOME=/opt/knlservices/Java/apache-maven-3.3.3
export MAVEN_HOME

PATH=$MAVEN_HOME/bin:$PATH

export PATH

MAVEN_LOCAL_REPO=/opt/knlservices/Java/local_maven_repo

export MAVEN_LOCAL_REPO

MAVEN_SETTINGS=/opt/knlservices/Java/local_maven_settings/settings.xml

export MAVEN_SETTINGS
